# systemd/update-proxmox-hosts.service
[Unit]
Description=Update /etc/hosts with current IP for Proxmox
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-proxmox-hosts.sh

# No [Install] section: this oneshot service is triggered by
# update-proxmox-hosts.timer rather than enabled on its own.
