# systemd/update-proxmox-hosts.service
[Unit]
Description=Update /etc/hosts with current IP for Proxmox
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-proxmox-hosts.sh

[Install]
WantedBy=multi-user.target
